# mongo_backend/crontab
20 18 * * * python3 /usr/src/app/stock_ingestion.py
